# provision/deamon.conf
[program:deamon-worker]
command= php /var/www/src/core/chatServer2.php
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/src/core/logs/worker.log